# launchd/news.briefer.dmarc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>news.briefer.dmarc</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/maxgoshay/code/briefernewsapp/scripts/dmarc_report.py</string>
  </array>

  <!-- 07:30: after the overnight scrape/synth, before the 09:00 alert digest, so
       any spoofing-alignment warning DMARC surfaces rides into that day's digest. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/maxgoshay/code/briefernewsapp/logs/dmarc-launchagent.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/maxgoshay/code/briefernewsapp/logs/dmarc-launchagent.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
