# launchd/news.briefer.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>news.briefer.backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/maxgoshay/code/briefernewsapp/scripts/backup_subscribers.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>23</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/maxgoshay/code/briefernewsapp/logs/backup-launchagent.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/maxgoshay/code/briefernewsapp/logs/backup-launchagent.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
