# launchd/news.briefer.spend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>news.briefer.spend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/maxgoshay/code/briefernewsapp/scripts/daily_memex.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/maxgoshay/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/maxgoshay/code/briefernewsapp</string>
  <key>StandardOutPath</key>
  <string>/Users/maxgoshay/code/briefernewsapp/logs/launchd-spend.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/maxgoshay/code/briefernewsapp/logs/launchd-spend.err</string>
</dict>
</plist>
