# launchd/news.briefer.healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>news.briefer.healthcheck</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/maxgoshay/code/briefernewsapp/scripts/healthcheck.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/maxgoshay/code/briefernewsapp/logs/healthcheck-launchagent.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/maxgoshay/code/briefernewsapp/logs/healthcheck-launchagent.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
